# 200-oc/opencode/systemd/opencode-copilot.service
[Unit]
Description=OpenCode Server (copilot)
After=network.target

[Service]
Type=simple
Environment=XDG_CONFIG_HOME=/home/jclee/.config/opencode-wt-copilot
Environment=PATH=/home/jclee/.nvm/versions/node/v22.22.0/bin:/home/jclee/.local/bin:/home/jclee/.opencode/bin:/usr/local/bin:/usr/bin:/bin
Environment=OPENCODE_SERVER_PASSWORD=opencode
ExecStart=/home/jclee/.opencode/bin/opencode serve --port 3003 --hostname 127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
