# 200-oc/config/systemd/opencode.socket
[Unit]
Description=OpenCode Socket Activation
Documentation=https://github.com/opencode-ai/opencode

[Socket]
ListenStream=/run/opencode/opencode.sock
SocketUser=jclee
SocketGroup=jclee
SocketMode=0660

[Install]
WantedBy=sockets.target
