# 200-oc/config/systemd/opencode-cleanup.service
[Unit]
Description=OpenCode Log Cleanup
Documentation=https://github.com/opencode-ai/opencode

[Service]
Type=oneshot
User=jclee
ExecStart=/home/jclee/proxmox/200-oc/scripts/cleanup-logs.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-cleanup
