# 200-oc/config/systemd/opencode.service
[Unit]
Description=OpenCode AI Coding Assistant
Documentation=https://github.com/opencode-ai/opencode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jclee
Group=jclee
WorkingDirectory=/home/jclee

# Environment
EnvironmentFile=-/home/jclee/.config/opencode/env
Environment=HOME=/home/jclee
Environment=XDG_CONFIG_HOME=/home/jclee/.config
Environment=XDG_DATA_HOME=/home/jclee/.local/share
Environment=PATH=/home/jclee/.opencode/bin:/usr/local/bin:/usr/bin:/bin

# Main process
ExecStart=/home/jclee/.opencode/bin/opencode serve --port 8090 --hostname 0.0.0.0 --print-logs
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=400%

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jclee/.config/opencode
ReadWritePaths=/home/jclee/.local/share/opencode
ReadWritePaths=/home/jclee/.opencode
ReadWritePaths=/home/jclee/.cache/opencode
ReadWritePaths=/home/jclee
ReadWritePaths=/tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode

[Install]
WantedBy=multi-user.target
