# 200-oc/config/systemd/opencode-anti.service
[Unit]
Description=OpenCode AI - Antigravity Profile
Documentation=https://github.com/opencode-ai/opencode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jclee
Group=jclee
WorkingDirectory=/home/jclee

# Environment
Environment=HOME=/home/jclee
Environment=PATH=/home/jclee/.opencode/bin:/usr/local/bin:/usr/bin:/bin

# Switch to anti profile before starting
ExecStartPre=/bin/ln -sf /home/jclee/.config/opencode/oh-my-opencode.json.anti /home/jclee/.config/opencode/oh-my-opencode.json

# Main process - port 3001
ExecStart=/home/jclee/.opencode/bin/opencode serve --port 3001 --hostname 127.0.0.1
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=400%

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jclee/.config/opencode
ReadWritePaths=/home/jclee/.local/share/opencode
ReadWritePaths=/home/jclee/.opencode
ReadWritePaths=/home/jclee/.cache/opencode
ReadWritePaths=/home/jclee
ReadWritePaths=/tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-anti

[Install]
WantedBy=multi-user.target
